# com.jamesno.retaildropmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for this LaunchAgent -->
    <key>Label</key>
    <string>com.jamesno.retaildropmonitor</string>

    <!-- Use the venv Python directly — no need to activate -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jamesno/Documents/retail-drop-monitor/venv/bin/python</string>
        <string>/Users/jamesno/Documents/retail-drop-monitor/monitor.py</string>
    </array>

    <!-- Run from the project folder so config.json and .env are found -->
    <key>WorkingDirectory</key>
    <string>/Users/jamesno/Documents/retail-drop-monitor</string>

    <!-- Start automatically when you log in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Log output to monitor.log (same file as nohup approach) -->
    <key>StandardOutPath</key>
    <string>/Users/jamesno/Documents/retail-drop-monitor/monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jamesno/Documents/retail-drop-monitor/monitor.log</string>
</dict>
</plist>
